# EPS Homebridge.indigoPlugin/Contents/Server Plugin/bin/hb/homebridge/1953136859/com.webdeck.homebridge.1953136859.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.webdeck.homebridge.1953136859</string>
	<key>WorkingDirectory</key>
	<string>/Users/Shared/Indigo Plugins/Homebridge/EPS Homebridge.indigoPlugin/Contents/Server Plugin/bin/hb/homebridge/1953136859</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/Shared/Indigo Plugins/Homebridge/EPS Homebridge.indigoPlugin/Contents/Server Plugin/bin/hb/homebridge/runhomebridge</string>
		<string>1953136859</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
